# package/gluon-core/files/etc/hotplug.d/net/10-gluon-apup-mesh
[ "$ACTION" = 'add' ] || exit 0

case "$INTERFACE" in
	apuppeer*)
		echo "Interface starts with ap"
		batctl if add "$INTERFACE"
		;;
	*)
		echo "Interface does not start with ap"
		;;
esac

